# templates/launchd/php.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>{{LABEL}}</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{FPM_EXEC}}</string>
    <string>{{VER}}</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HARBOR_ROOT</key><string>{{HARBOR_ROOT}}</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>{{LOG}}</string>
  <key>StandardErrorPath</key><string>{{LOG}}</string>
  <key>SoftResourceLimits</key>
  <dict><key>NumberOfFiles</key><integer>10240</integer></dict>
</dict>
</plist>
